# ghoma2mqtt/rootfs/etc/services.d/example/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start the example service
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

# Add your code here

# Declare variables
#declare message

## Get the 'message' key from the user config options.
#message=$(bashio::config 'message')

## Print the message the user supplied, defaults to "Hello World..."
#bashio::log.info "${message:="Hello World..."}"

## Run your program
export MQTT_USER="$(bashio::services 'mqtt' 'username')"
export MQTT_PASSWORD="$(bashio::services 'mqtt' 'password')"
if bashio::var.true "$(bashio::services 'mqtt' 'ssl')"; then
    export MQTT_SERVER="mqtts://$(bashio::services 'mqtt' 'host'):$(bashio::services 'mqtt' 'port')"
else
    export MQTT_SERVER="mqtt://$(bashio::services 'mqtt' 'host'):$(bashio::services 'mqtt' 'port')"
fi
exec /usr/bin/node /root/ghoma-mqtt-connect.js
